# templates/default/kafka.service.erb
[Unit]
Description=Apache Kafka server (broker)
Documentation=http://kafka.apache.org/documentation.html
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
Type=simple
PIDFile=/var/run/kafka.pid
EnvironmentFile=/etc/default/kafka
User=${KAFKA_USER}
Group=${KAFKA_USER}
ExecStart=${KAFKA_BIN}/kafka-server-start.sh ${KAFKA_CONFIG}/server.properties
Restart=on-failure
SyslogIdentifier=kafka

[Install]
WantedBy=multi-user.target
